# deploy/jprq.service
[Unit]
Description=JPRQ self-hosted tunnel server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
DynamicUser=yes
EnvironmentFile=/etc/jprq/jprq.env
ExecStart=/usr/local/lib/jprq/jprq-server
Restart=on-failure
RestartSec=3s
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
CapabilityBoundingSet=
SystemCallArchitectures=native
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
